# v1/etc/init.d/rc.local
#!/bin/sh

#/root/report/report_loop.sh &

